# etc/systemd/system/mnt-blobfuse.automount
[Unit]
Description=Automount blobfuse
ConditionPathExists=/mnt/blobfuse

[Automount]
Where=/mnt/blobfuse
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
